# patches/gtk+2/profile.d/gtk+.sh.mirrorlist
<!DOCTYPE html PUBLIC "-//W3C//DTD XHTML 1.0 Transitional//EN"
"http://www.w3.org/TR/xhtml1/DTD/xhtml1-transitional.dtd">
<html xmlns="http://www.w3.org/1999/xhtml">
<head>
  <title>Mirror List</title>
  <link type="text/css" rel="stylesheet" href="/css/mblist.css" />
</head>

<body>
<div id="mirrorbrain-details">
  <h2>Mirrors for <a href="http://mirrors.slackware.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirrors.slackware.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a></h2>
  <br/>
  <address>Powered by <a href="http://mirrorbrain.org/">MirrorBrain</a> on <a href="http://slackware.com/">Slackware</a></address>
  <ul>
  <li>Size:  44  (44 bytes)</li>
  <li>Last modified: Thu, 09 Dec 2010 19:13:01 GMT (Unix time: 1291921981)</li>
  <li><a href="http://mirrors.slackware.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh.sha256">SHA-256 Hash</a>: <tt>443e94e68d09d68255c0b595e3acba45a2d9efdc5b71d7324b77c1776509ec88</tt> </li>
  <li><a href="http://mirrors.slackware.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh.sha1">SHA-1 Hash</a>: <tt>6fa4685acdc7ac29be4994259a7989aac71f1327</tt> </li>
  <li><a href="http://mirrors.slackware.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh.md5">MD5 Hash</a>: <tt>2346d87621891129b13bc18d3287b930</tt> </li>
  </ul>
  <a href="http://mirrors.slackware.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh.meta4">http://mirrors.slackware.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh.meta4</a> (IETF Metalink)  <br/>
  </blockquote>
  <p>List of best mirrors for IP address 65.94.103.145, located at 43.860298,-79.043404 in Canada (CA): &nbsp;
     <a href="http://maps.google.com/maps/api/staticmap?size=640x640&amp;maptype=terrain&amp;visible&amp;sensor=false&amp;markers=size:mid|color:red|43.860298,-79.043404&amp;markers=size:normal|color:yellow|label:1|43.466000,-80.501999&amp;markers=size:normal|color:yellow|label:2|44.650002,-63.599998&amp;markers=size:normal|color:yellow|label:3|42.278000,-83.740997&amp;markers=size:normal|color:yellow|label:4|39.563999,-75.597000&amp;markers=size:normal|color:yellow|label:5|42.726002,-84.636002&amp;markers=size:normal|color:yellow|label:6|43.305000,-70.975998&amp;markers=size:normal|color:yellow|label:7|41.875000,-87.650002&amp;markers=size:normal|color:yellow|label:8|47.153999,-88.647003&amp;markers=size:normal|color:yellow|label:9|43.074001,-89.527000">
     <img src="http://maps.google.com/maps/api/staticmap?size=50x50&amp;maptype=terrain&amp;visible&amp;sensor=false&amp;markers=size:mid|color:red|43.860298,-79.043404&amp;markers=size:normal|color:yellow|label:1|43.466000,-80.501999&amp;markers=size:normal|color:yellow|label:2|44.650002,-63.599998&amp;markers=size:normal|color:yellow|label:3|42.278000,-83.740997&amp;markers=size:normal|color:yellow|label:4|39.563999,-75.597000&amp;markers=size:normal|color:yellow|label:5|42.726002,-84.636002&amp;markers=size:normal|color:yellow|label:6|43.305000,-70.975998&amp;markers=size:normal|color:yellow|label:7|41.875000,-87.650002&amp;markers=size:normal|color:yellow|label:8|47.153999,-88.647003&amp;markers=size:normal|color:yellow|label:9|43.074001,-89.527000" width="50" height="50" alt="map showing the closest mirrors"/></a> </p>

  <h3>Found 2 mirrors which handle this country (CA):</h3>
  <ul>
    <li><a href="http://mirror.csclub.uwaterloo.ca/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.csclub.uwaterloo.ca/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (ca, prio 100)</li>
    <li><a href="http://mirror.its.dal.ca/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.its.dal.ca/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (ca, prio 100)</li>
  </ul>

  <h3>Found 21 mirrors in other countries, but same continent (NA):</h3>
  <ul>
    <li><a href="http://mirror.nexcess.net/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.nexcess.net/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://mirror.lug.udel.edu/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.lug.udel.edu/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://mirrors.kingrst.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirrors.kingrst.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://mirror.metrocast.net/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.metrocast.net/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://slackblog.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://slackblog.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://lug.mtu.edu/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://lug.mtu.edu/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://slackware.mirrors.tds.net/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://slackware.mirrors.tds.net/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://ftp.gtlib.gatech.edu/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.gtlib.gatech.edu/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://slackbuilds.org/mirror/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://slackbuilds.org/mirror/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://dfw.mirror.rackspace.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://dfw.mirror.rackspace.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://slackware.mirrorcatalogs.com/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://slackware.mirrorcatalogs.com/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://slackware.cs.utah.edu/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://slackware.cs.utah.edu/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://mirrors.xmission.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirrors.xmission.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://taper.alienbase.nl/mirrors/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://taper.alienbase.nl/mirrors/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://mirrors.usc.edu/pub/linux/distributions/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirrors.usc.edu/pub/linux/distributions/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://mirrors4.kernel.org/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirrors4.kernel.org/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://mirrors1.kernel.org/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirrors1.kernel.org/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://mirrors2.kernel.org/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirrors2.kernel.org/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://mirrors3.kernel.org/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirrors3.kernel.org/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
    <li><a href="http://ftp.slackware.com/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.slackware.com/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 1)</li>
    <li><a href="http://carroll.cac.psu.edu/pub/linux/distributions/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://carroll.cac.psu.edu/pub/linux/distributions/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (us, prio 100)</li>
  </ul>

   <h3>Found 44 mirrors in other parts of the world:</h3>
  <ul>
    <li><a href="http://mirror.cedia.org.ec/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.cedia.org.ec/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (ec, prio 100)</li>
    <li><a href="http://mirror.ic.ufmt.br/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.ic.ufmt.br/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (br, prio 100)</li>
    <li><a href="http://ftp.heanet.ie/mirrors/ftp.slackware.com/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.heanet.ie/mirrors/ftp.slackware.com/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (ie, prio 100)</li>
    <li><a href="http://linorg.usp.br/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://linorg.usp.br/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (br, prio 100)</li>
    <li><a href="http://www.mirrorservice.org/sites/ftp.slackware.com/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://www.mirrorservice.org/sites/ftp.slackware.com/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (gb, prio 100)</li>
    <li><a href="http://slackware.org.uk/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://slackware.org.uk/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (gb, prio 100)</li>
    <li><a href="http://mirror.bytemark.co.uk/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.bytemark.co.uk/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (gb, prio 100)</li>
    <li><a href="http://lon.mirror.rackspace.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://lon.mirror.rackspace.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (gb, prio 100)</li>
    <li><a href="http://slackware.mirrors.ovh.net/ftp.slackware.com/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://slackware.mirrors.ovh.net/ftp.slackware.com/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (fr, prio 100)</li>
    <li><a href="http://ftp.belnet.be/ftp.slackware.com/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.belnet.be/ftp.slackware.com/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (be, prio 100)</li>
    <li><a href="http://mirror.nl.leaseweb.net/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.nl.leaseweb.net/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (nl, prio 100)</li>
    <li><a href="http://slackware.bokxing-it.nl/mirror/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://slackware.bokxing-it.nl/mirror/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (nl, prio 100)</li>
    <li><a href="http://mirror.netcologne.de/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.netcologne.de/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (de, prio 100)</li>
    <li><a href="http://ftp5.gwdg.de/pub/linux/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp5.gwdg.de/pub/linux/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (de, prio 100)</li>
    <li><a href="http://ftp.slackware.no/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.slackware.no/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (no, prio 100)</li>
    <li><a href="http://www.slackware.at/data/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://www.slackware.at/data/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (at, prio 100)</li>
    <li><a href="http://mirror.slackware.hr/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.slackware.hr/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (hr, prio 100)</li>
    <li><a href="http://mirror.oss.maxcdn.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.oss.maxcdn.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (cz, prio 100)</li>
    <li><a href="http://mirror.wheel.sk/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.wheel.sk/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (sk, prio 100)</li>
    <li><a href="http://piotrkosoft.net/pub/mirrors/ftp.slackware.com/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://piotrkosoft.net/pub/mirrors/ftp.slackware.com/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (pl, prio 100)</li>
    <li><a href="http://sunsite.icm.edu.pl/pub/Linux/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://sunsite.icm.edu.pl/pub/Linux/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (pl, prio 100)</li>
    <li><a href="http://slackware.blizoo.mk/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://slackware.blizoo.mk/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (mk, prio 100)</li>
    <li><a href="http://atviras.lt/veidrodziai/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://atviras.lt/veidrodziai/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (lt, prio 100)</li>
    <li><a href="http://slackware.laukas.lt/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://slackware.laukas.lt/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (lt, prio 100)</li>
    <li><a href="http://x.ultranet.biz/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://x.ultranet.biz/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (bg, prio 100)</li>
    <li><a href="http://ftp.cc.uoc.gr/mirrors/linux/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.cc.uoc.gr/mirrors/linux/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (gr, prio 100)</li>
    <li><a href="http://slackware.telecoms.bg/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://slackware.telecoms.bg/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (bg, prio 100)</li>
    <li><a href="http://mirrors.nix.org.ua/linux/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirrors.nix.org.ua/linux/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (ua, prio 100)</li>
    <li><a href="http://ftp.linux.org.tr/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.linux.org.tr/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (tr, prio 100)</li>
    <li><a href="http://mirror.yandex.ru/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.yandex.ru/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (ru, prio 100)</li>
    <li><a href="http://ftp.wa.co.za/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.wa.co.za/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (za, prio 100)</li>
    <li><a href="http://ftp.is.co.za/mirror/ftp.slackware.com/pub/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.is.co.za/mirror/ftp.slackware.com/pub/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (za, prio 100)</li>
    <li><a href="http://slackware.tsu.ru/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://slackware.tsu.ru/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (ru, prio 100)</li>
    <li><a href="http://syd.mirror.rackspace.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://syd.mirror.rackspace.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (au, prio 100)</li>
    <li><a href="http://hkg.mirror.rackspace.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://hkg.mirror.rackspace.com/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (hk, prio 100)</li>
    <li><a href="http://mirror.bjtu.edu.cn/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.bjtu.edu.cn/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (cn, prio 100)</li>
    <li><a href="http://mirror.pregi.net/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.pregi.net/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (ph, prio 100)</li>
    <li><a href="http://ftp.nara.wide.ad.jp/pub/Linux/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.nara.wide.ad.jp/pub/Linux/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (jp, prio 100)</li>
    <li><a href="http://mirror.aarnet.edu.au/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.aarnet.edu.au/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (au, prio 100)</li>
    <li><a href="http://mirror.internode.on.net/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.internode.on.net/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (au, prio 100)</li>
    <li><a href="http://ftp.swin.edu.au/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.swin.edu.au/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (au, prio 100)</li>
    <li><a href="http://mirror.lagoon.nc/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://mirror.lagoon.nc/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (nc, prio 100)</li>
    <li><a href="http://ftp.sunet.se/pub/os/Linux/distributions/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.sunet.se/pub/os/Linux/distributions/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (se, prio 100)</li>
    <li><a href="http://ftp.slackware.pl/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh">http://ftp.slackware.pl/pub/slackware/slackware-14.1/source/l/gtk+2/profile.d/gtk+.sh</a> (pl, prio 100)</li>
  </ul>
</div> <!-- mirrorbrain-details -->
<!--
<html>
<a href="http://mirrorbrain.org/">MirrorBrain</a>
powered by <a href="http://slackware.com/">Slackware</a>!
</html>
-->
